# .config/systemd/user/checkmail.timer
[Unit]
Description=Check email

[Timer]
OnStartupSec=1m
OnUnitActiveSec=15m
Persistent=true

[Install]
WantedBy=timers.target
