# .config/systemd/user/checkupdates.service
[Unit]
Description=Arch Linux updates
OnFailure=notify-send@%i.service
Wants=NetworkManager-wait-online.service
After=NetworkManager-wait-online.service

[Service]
Type=oneshot
ExecStart=/usr/bin/checkupdates
ExecStartPost=-/usr/bin/polybar-msg hook checkupdates 1
SuccessExitStatus=SUCCESS 2
