# .config/systemd/user/checkupdates.service
[Unit]
Description=Arch Linux updates

[Service]
Type=oneshot
ExecStart=/usr/bin/checkupdates
SuccessExitStatus=SUCCESS 2
