# .config/systemd/user/context.service
[Unit]
Description=Update ConTeXt standalone

[Service]
Type=oneshot
ExecStart=/bin/bash /opt/context/update.sh
