# .config/systemd/user/context-lmtx.service
[Unit]
Description=Update ConTeXt LMTX
OnFailure=notify-send@%i.service

[Service]
Type=oneshot
ExecStart=/bin/bash /opt/context-lmtx/update.sh
