# .config/systemd/user/context.timer
[Unit]
Description=Update ConTeXt standalone

[Timer]
OnCalendar=daily
Persistent=true

[Install]
WantedBy=timers.target
