# .config/systemd/user/context-lmtx.timer
[Unit]
Description=Update ConTeXt LMTX

[Timer]
OnBootSec=10min
OnUnitActiveSec=1d
Persistent=true

[Install]
WantedBy=timers.target
