# .config/systemd/user/context-lmtx.service
[Unit]
Description=Update ConTeXt LMTX

[Service]
Type=oneshot
ExecStart=/bin/bash /opt/context-lmtx/update.sh
