# .config/systemd/user/mbsync@.timer
[Unit]
Description=mbsync timer for account %i

[Timer]
OnBootSec=1m
OnUnitInactiveSec=15m

[Install]
WantedBy=default.target
