# .config/systemd/user/newsboat.timer
[Unit]
Description=Update news feeds

[Timer]
OnBootSec=5min
OnUnitActiveSec=5min
Persistent=true

[Install]
WantedBy=timers.target
