# .config/systemd/user/newsboat.service
[Unit]
Description=Update news feeds

[Service]
Type=oneshot
ExecStart=/usr/bin/newsboat -x reload
ExecStartPost=-/usr/bin/polybar-msg hook newsboat 1
