# .config/systemd/user/newsboat.timer
[Unit]
Description=Update news feeds

[Timer]
OnBootSec=5m
OnUnitActiveSec=1h
Persistent=true

[Install]
WantedBy=timers.target
