# .config/systemd/user/syncthing-gui.socket
[Socket]
ListenStream=8384

[Install]
WantedBy=sockets.target
