# .config/systemd/user/texlive@.timer
[Unit]
Description=Update TeX Live %i

[Timer]
OnBootSec=5min
OnUnitActiveSec=1d
Persistent=true

[Install]
WantedBy=timers.target
